# app-misc/mcp-log-server/files/mcp-log-server.service
[Unit]
Description=MCP Log Server (FastAPI)
After=network.target

[Service]
User=mcp-log
Group=mcp-log
WorkingDirectory=/usr/libexec/mcp-log-server
ExecStart=/usr/libexec/mcp-log-server/venv/bin/python \
    /usr/libexec/mcp-log-server/log_mcp.py
Restart=always

[Install]
WantedBy=multi-user.target