# app-misc/mcp-log-server/files/mcp-log-server.service
[Unit]
Description=MCP Log Server (FastAPI)
After=network.target

[Service]
User=mcp-log
Group=mcp-log
WorkingDirectory=/usr/libexec/mcp-log-server
ExecStart=/usr/libexec/mcp-log-server/venv/bin/uvicorn \
    log_mcp:app \
    --host 0.0.0.0 \
    --port 8000
Restart=always

[Install]
WantedBy=multi-user.target