# app-misc/mcp-rag/files/mcp-rag.service
[Unit]
Description=MCP RAG Server (FAISS + SentenceTransformers)
After=network.target

[Service]
User=mcp-rag
Group=mcp-rag
Type=simple
ExecStart=/usr/bin/uvicorn rag_server:app --host 127.0.0.1 --port 8003
WorkingDirectory=/usr/libexec/mcp-rag
Restart=on-failure
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=default.target
