# app-misc/mcp-fs-server/files/mcp-fs-server.service
[Unit]
Description=MCP Filesystem Server for /home/oli/together
After=network.target

[Service]
User=mcp-fs
Group=mcp-fs
Type=simple
ExecStart=/usr/bin/uvicorn filesystem_mcp:app --host 127.0.0.1 --port 8001
WorkingDirectory=/usr/libexec/mcp-fs-server
Restart=on-failure
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=default.target
